# pi_config_files/pifinder_splash.service
[Unit]
Description=PiFinderSplash

[Service]
WorkingDirectory=/home/pifinder/PiFinder5/python
ExecStart=/usr/local/lib/python3.10 -m PiFinder.splash

[Install]
WantedBy=basic.target
